# RUFilm.service
[Service]
WorkingDirectory=/home/RUFilm
User=root
ExecStart=/usr/bin/python3 main.py

[Install]
WantedBy=multi-user.target